# deb/src/DEBIAN/postinst
#!/bin/sh

# Copyright:   (c) 2018, Vladimir Vons, UA
# Author:      Vladimir Vons <VladVons@gmail.com>
# Created:     01.11.2018
# License:     GNU, see LICENSE for more details
# Description:

# Exit on error
set -e

#. /usr/share/debconf/confmodule
Name="py-relay"


Configure()
{
    # Setup logfile permissions:
    LogDir=/var/log/$Name
    mkdir -p $LogDir

    touch $LogDir/$Name.log
    chmod 640 $LogDir/$Name.log

    if [ -x "/etc/init.d/$Name" ]; then
        update-rc.d $Name defaults
        update-rc.d $Name enable
    fi
}


Debug()
{
  aMode="$1";

  Msg="postinst: $aMode"
  echo $Msg
  #echo $Msg >> /$Name.log
  #printenv > /$Name.log2
}

#Debug $1

case "$1" in
  configure)
    Configure
    ;;
  abort-upgrade|abort-remove|abort-deconfigure)
    ;;
  *)
    echo "postinst called with unknown argument '$1'" >&2
    exit 1
    ;;
esac

exit 0
